# scripts/gowebserver.service
[Unit]
Description=Go Web Server
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=gowebserver
Group=gowebserver
WorkingDirectory=/opt/gowebserver
ExecStart=/opt/gowebserver/bin/server
Restart=always
RestartSec=10
KillMode=mixed
TimeoutStopSec=5
EnvironmentFile=/opt/gowebserver/.env

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/gowebserver/logs
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target